# telegram-assistant/mac-indexer/com.alexander.file-indexer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alexander.file-indexer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/user/file-indexer/indexer.py</string>
    </array>

    <!-- Запуск сразу при загрузке -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Повтор каждые 2 часа (7200 секунд) -->
    <key>StartInterval</key>
    <integer>7200</integer>

    <key>StandardOutPath</key>
    <string>/Users/user/file-indexer/indexer.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/user/file-indexer/indexer.error.log</string>

    <!-- Запускать во время Power Nap (Mac спит с питанием) -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Не перезапускать при сбое (избежать бесконечного цикла) -->
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
